# deploy/systemd/getbible-search.service.d/limits.conf
# Starting whole-unit limits for the full-corpus Search service.
# Calibrate against measured workers, translations, and index variants.
[Unit]
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
MemoryAccounting=true
TasksAccounting=true
MemoryHigh=6G
MemoryMax=8G
MemorySwapMax=0
TasksMax=512
LimitNOFILE=65536
OOMPolicy=stop
Restart=on-failure
RestartSec=5s
TimeoutStopSec=45s
KillMode=mixed
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

ReadOnlyPaths=/srv/getbible/mirror
ReadWritePaths=/var/cache/getbible/search
