# deploy/systemd/getbible-query.service.d/limits.conf
# Starting whole-unit limits for the lightweight Query service.
# Calibrate against measured worker count and chapter-cache residency.
[Unit]
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
MemoryAccounting=true
TasksAccounting=true
MemoryHigh=1536M
MemoryMax=2G
MemorySwapMax=0
TasksMax=256
LimitNOFILE=65536
OOMPolicy=stop
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
KillMode=mixed
UMask=0027

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectClock=true
ProtectHostname=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
RestrictRealtime=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

ReadOnlyPaths=/srv/getbible/mirror
ReadWritePaths=/var/cache/getbible/query
