# deploy/systemd/getbible-query.service.d/limits.conf
[Service]
MemoryAccounting=yes
MemoryHigh=384M
MemoryMax=512M
MemorySwapMax=0
CPUAccounting=yes
CPUQuota=200%
TasksAccounting=yes
TasksMax=64
LimitNOFILE=4096
OOMPolicy=stop
TimeoutStartSec=90
TimeoutStopSec=30
